# systemd/speaker-watchdog.service
[Unit]
Description=Speaker Watchdog Service for Home Assistant
Documentation=https://github.com/danuser2018/speaker-watchdog
After=default.target sound.target

[Service]
Type=simple
# La ruta al WorkingDirectory, EnvironmentFile y al ejecutable se establecen
# durante la instalación mediante el script scripts/install.sh.
# No edites este archivo directamente; el archivo real instalado en
# ~/.config/systemd/user/ contiene las rutas absolutas.
WorkingDirectory=__SPEAKER_WATCHDOG_DIR__
EnvironmentFile=__SPEAKER_WATCHDOG_ENV_FILE__
ExecStart=__SPEAKER_WATCHDOG_VENV_PYTHON__ src/main.py
Restart=always
RestartSec=3
Environment=XDG_RUNTIME_DIR=/run/user/__USER_UID__
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
